# sandbox1/docker/php/worker.conf
; ==============================================================================
; Supervisor 基本設定 (UNIXソケット)
; ==============================================================================
[unix_http_server]
; 管理用ソケットファイルのパスを指定 (supervisorctlとの通信インターフェースとして利用するため)
file=/var/run/supervisor.sock
; アクセス権限を制限 (rootユーザーのみに操作権限を限定しセキュリティを確保するため)
chmod=0700

; ==============================================================================
; Supervisord デーモン本体の設定
; ==============================================================================
[supervisord]
; フォアグラウンドで実行 (コンテナのメインプロセスとしてPID 1を維持し停止を防ぐため)
nodaemon=true
; 実行ユーザーをrootに設定 (サブプロセスごとに異なるユーザーへ切り替える権限を持たせるため)
user=root
; ログを標準出力へ転送 (docker logsコマンドで全プロセスの挙動を一括監視するため)
logfile=/dev/stdout
; ローテーションを無効化 (ログの肥大化管理をDockerエンジン側の機能に委ねるため)
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
; 環境変数を子プロセスへ継承 (全プロセスで一貫したタイムゾーンを適用させるため)
environment=TZ="Asia/Tokyo",APP_TIMEZONE="Asia/Tokyo"

; ==============================================================================
; RPCインターフェース & クライアント設定
; ==============================================================================
[rpcinterface:supervisor]
; RPCインターフェースを有効化 (外部ツールやCLIからのプロセス制御を可能にするため)
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
; 接続先ソケットを指定 (CLIツールがデーモンと通信するためのエンドポイントを確立するため)
serverurl=unix:///var/run/supervisor.sock

; ==============================================================================
; プログラム設定: Laravel スケジューラー
; ==============================================================================
[program:laravel-scheduler]
; 1分間隔のタスク実行を維持 (cronを使わずコンテナ内で完結した定期実行を実現するため)
command=bash -c "exec /usr/local/bin/php -d output_buffering=0 /var/www/html/artisan schedule:work 2>&1 | stdbuf -oL sed 's/^/[SCHEDULER] /'"
directory=/var/www/html
; 自動起動を有効化 (コンテナ起動と同時にタスクスケジュールを開始させるため)
autostart=true
; 異常終了時の再起動を有効化 (一時的なエラーによるタスク停止を自動復旧させるため)
autorestart=true
; 実行ユーザーをwww-dataに指定 (アプリケーションが生成するファイルの所有権を一貫させるため)
user=www-data
; 停止時の猶予時間を設定 (実行中のタスクが中途半端に終了するのを防ぐため)
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; 起動優先度を高く設定 (Webサーバーより先に基盤プロセスを立ち上げるため)
priority=1

; ==============================================================================
; プログラム設定: Apache2 (Webサーバー)
; ==============================================================================
[program:apache2]
; Apacheをフォアグラウンド実行 (Webリクエストの受付をSupervisorの管理下で継続するため)
command=bash -c "exec apache2-foreground 2>&1 | stdbuf -oL sed 's/^/[APACHE] /'"
directory=/var/www/html
autostart=true
autorestart=true
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; 中間の優先度を設定 (スケジューラーの準備が整った後にリクエスト受付を開始するため)
priority=10

; ==============================================================================
; プログラム設定: Laravel キューワーカー
; ==============================================================================
[program:laravel-worker]
; キュー処理を常駐実行し長時間ジョブのメモリリークを防止 (--max-timeで定期再起動させるため)
command=bash -c "exec /usr/local/bin/php -d output_buffering=0 /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-jobs=1000 --max-time=3600 --queue=high,default,low 2>&1 | stdbuf -oL sed 's/^/[WORKER] /'"
directory=/var/www/html
autostart=true
autorestart=true
user=www-data
; 停止猶予を1時間に設定 (長時間実行されるジョブの安全な終了を保証するため)
stopwaitsecs=3600
; プロセスグループ全体を停止 (ゾンビプロセスの発生を抑制しリソースをクリーンに保つため)
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; 起動優先度を低く設定 (DBやRedisへの接続準備が完全に整ってからジョブ処理を開始するため)
priority=999
