# sandbox1/docker/php/worker.conf
# ==============================================================================
# 1. Supervisor 基本設定 (UNIXソケット)
# supervisorctl コマンド（状態確認・再起動）を使用するために必須の設定です
# ==============================================================================
[unix_http_server]
file=/var/run/supervisor.sock   ; supervisorctlが通信に使うソケットファイル
chmod=0700                      ; rootユーザーのみアクセス許可

# ==============================================================================
# 2. Supervisord デーモン本体の設定
# ==============================================================================
[supervisord]
nodaemon=true                    ; Dockerコンテナを終了させないためにフォアグラウンド実行
user=root                        ; デーモン自体はrootで起動（各プログラムでユーザー切替）
logfile=/dev/fd/1                ; 全体のログを標準出力に流し、docker logs で確認可能にする
logfile_maxbytes=0               ; ログローテーションを無効化（Docker側で管理するため）
pidfile=/var/run/supervisord.pid ; プロセスIDファイルの場所
environment=TZ="Asia/Tokyo",APP_TIMEZONE="Asia/Tokyo" ; 子プロセスにタイムゾーンを継承

# ==============================================================================
# 3. RPCインターフェース & クライアント設定 (supervisorctlの動作に必要)
# ==============================================================================
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; 上記 [unix_http_server] のパスと一致させる

# ==============================================================================
# 4. プログラム設定: Laravel スケジューラー
# 1分ごとにタスクがあるかチェックする常駐プロセス
# ==============================================================================
[program:laravel-scheduler]
# -d output_buffering=0: PHPのバッファリングを無効化し、ログを即座に出力させる
# ログ検索用に識別子'[SCHEDULER]'を付与
command=bash -c "exec /usr/local/bin/php -d output_buffering=0 /var/www/html/artisan schedule:work 2>&1 | stdbuf -oL sed 's/^/[SCHEDULER] /'"
directory=/var/www/html
autostart=true                  ; Supervisor起動時に自動開始
autorestart=true                ; 落ちても自動で再起動
user=www-data                   ; Webサーバー実行ユーザーで動かす
stdout_logfile=/dev/fd/1        ; ログを標準出力に統合
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1                      ; 起動優先度（低）

# ==============================================================================
# 5. プログラム設定: Apache2 (Webサーバー)
# Webリクエストを処理するメインサーバー
# ==============================================================================
[program:apache2]
# Apacheをフォアグラウンドで実行（これがないとコンテナが終了してしまう）
# ログ検索用に識別子'[APACHE]'を付与
command=bash -c "exec /usr/sbin/apache2ctl -D FOREGROUND 2>&1 | stdbuf -oL sed 's/^/[APACHE] /'"
directory=/var/www/html
autostart=true
autorestart=true
user=root
stdout_logfile=/dev/fd/1        ; アクセスログ等を標準出力に統合
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10                     ; スケジューラーの後に起動

# ==============================================================================
# 6. プログラム設定: Laravel キューワーカー
# jobsテーブルを監視し、非同期処理を実行するプロセス
# ==============================================================================
[program:laravel-worker]
# --queue: 優先度順。 --sleep: ジョブ不在時の待機秒数。 --tries: 失敗時のリトライ回数。
# ログ検索用に識別子'[WORKER]'を付与
command=bash -c "exec /usr/local/bin/php -d output_buffering=0 /var/www/html/artisan queue:work --sleep=3 --tries=3 --queue=high,default,low 2>&1 | stdbuf -oL sed 's/^/[WORKER] /'"
directory=/var/www/html
autostart=true
autorestart=true
user=www-data
# 停止時に実行中のジョブを保護する設定 (Graceful Shutdown)
stopwaitsecs=3600               ; ジョブ完了まで最大1時間待つ
stopasgroup=true                ; 子プロセスも含めて停止
killasgroup=true                ; 強制終了時もグループ全体を対象にする
stdout_logfile=/dev/fd/1        ; 実行ログを標準出力に統合
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
# 出力を溜め込まず即座にイベントとして発行する
stdout_events_enabled=true
stderr_events_enabled=true
priority=999                    ; 最後に起動させる
