# sandbox1/docker/php/worker.conf
;==============================================================================
; Supervisor 基本設定 (UNIXソケット)
;==============================================================================
[unix_http_server]
; supervisorctl との通信インターフェースとして使用するソケットファイル
file=/var/run/supervisor.sock
; root ユーザーのみに操作権限を限定しセキュリティを確保
chmod=0700

;==============================================================================
; Supervisord デーモン本体の設定
;==============================================================================
[supervisord]
; コンテナの PID 1 として維持するためフォアグラウンドで実行
nodaemon=true
; サブプロセスごとに異なるユーザーへ切り替える権限を持たせるため root で実行
user=root
; docker logs コマンドで全プロセスの挙動を一括監視できるよう標準出力へ転送
logfile=/dev/stdout
; ログのローテーション管理を Docker エンジン側に委ねるためローテーションを無効化
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
; 全プロセスで一貫したタイムゾーンを適用させるため子プロセスへ環境変数を継承
environment=TZ="Asia/Tokyo",APP_TIMEZONE="Asia/Tokyo"

;==============================================================================
; RPCインターフェース & クライアント設定
;==============================================================================
[rpcinterface:supervisor]
; 外部ツールや CLI からのプロセス制御を可能にするため RPC インターフェースを有効化
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
; CLI ツールがデーモンと通信するためのエンドポイント
serverurl=unix:///var/run/supervisor.sock

;==============================================================================
; プログラム設定: Laravel スケジューラー
;==============================================================================
[program:laravel-scheduler]
; cron を使わずコンテナ内で完結した定期実行を実現するため schedule:work を常駐起動
; output_buffering=0 でバッファリングなしの即時ログ出力を保証
command=bash -c "exec /usr/local/bin/php -d output_buffering=0 /var/www/html/artisan schedule:work 2>&1 | stdbuf -oL sed 's/^/[SCHEDULER] /'"
directory=/var/www/html
autostart=true
; 一時的なエラーによるタスク停止を自動復旧
autorestart=true
; アプリケーションが生成するファイルの所有権を一貫させるため www-data で実行
user=www-data
; 実行中タスクの中途終了を防ぐため停止猶予を設定
stopwaitsecs=60
; ゾンビプロセスの発生を抑制するためプロセスグループ全体を停止対象に含める
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Web サーバーより先に起動させる基盤プロセスのため優先度を高く設定
priority=1

;==============================================================================
; プログラム設定: Apache2 (Webサーバー)
;==============================================================================
[program:apache2]
; Supervisor の管理下で Web リクエストの受付を継続するためフォアグラウンド実行
command=bash -c "exec apache2-foreground 2>&1 | stdbuf -oL sed 's/^/[APACHE] /'"
directory=/var/www/html
autostart=true
autorestart=true
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; スケジューラー準備後にリクエスト受付を開始させる中間優先度
priority=10

;==============================================================================
; プログラム設定: Laravel キューワーカー
;==============================================================================
[program:laravel-worker]
; --max-time で定期再起動させてメモリリークを防止し長時間ジョブを安定稼働
command=bash -c "exec /usr/local/bin/php -d output_buffering=0 /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-jobs=1000 --max-time=3600 --queue=high,default,low 2>&1 | stdbuf -oL sed 's/^/[WORKER] /'"
directory=/var/www/html
autostart=true
autorestart=true
user=www-data
; 長時間実行ジョブの安全な終了を保証するため停止猶予を 1 時間に設定
stopwaitsecs=3600
; ゾンビプロセスの発生を抑制するためプロセスグループ全体を停止対象に含める
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; DB や Redis の準備が整ってからジョブ処理を開始させるため優先度を最低に設定
priority=999
