# sandbox1/docker/php/worker.conf
[program:laravel-worker]
# --- 実行コマンド ---
# . /etc/environment; を入れることで、Docker の環境変数 (DBパスワード等) を PHP 伝える
# --queue=high,default,low (キー名)は主にこの3種類を使う
# PHP 8.5環境では --max-jobs や --max-time を入れることで、メモリ肥大化を防ぐ
command=sh -c ". /etc/environment; /usr/local/bin/php /var/www/html/base/artisan queue:work --sleep=3 --tries=3 --queue=high,default,low --max-jobs=1000 --max-time=3600"
directory=/var/www/html/base

# --- 起動・再起動設定 ---
autostart=true
autorestart=true
user=www-data

# --- 停止時の挙動 (Graceful Shutdown) ---
# コンテナ停止時に「処理中のジョブ」を強制終了させず、最大1時間待つ設定
# これがないと、DB整合性が壊れるリスクがある
stopwaitsecs=3600
stopasgroup=true
killasgroup=true

# --- ログ設定 (Docker最適化) ---
# 標準出力へ流すことで `docker logs supervisor` で確認可能
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# 優先順位で他のプロセスより後に起動するように設定
priority=999
